# config/supervisords_gunicorn.ini
[program:gunicorn]                                                           # gunicorn为进程的名字
user=qiushenjie	                                                             # 操作的用户
directory=/home/ubuntu/qiushenjie/FlaskModel                              # 项目目录，
command=/home/ubuntu/anaconda3/envs/qsj/bin/gunicorn -c gun.py manage:app    # 启动flask服务的命令，gunicorn命令路径通过whereis gunicorn查看，不然在supervisorctl start gunicorn时会报错
startsecs=5                                                                  # 启动5秒后没有异常退出，视作正常启动
autostart=true                                                               # 在 supervisord 启动时自动启动
autorestart=true                                                             # 程序异常退出后重启
redirect_stderr=true                                                         # 将错误信息重定向至stdout日志
stdout_logfile=/home/ubuntu/qiushenjie/FlaskModel/logs/gunicorn.pid    # 进程日志